# devops_tasks/systemd/backup.service
[Unit]
Description=Backup /etc directory

[Service]
ExecStart=/usr/local/bin/backup.sh
